# roles/restic-postgres/templates/restic-postgres-backup.timer.j2
[Unit]
Description=Run restic-postgres-backup.service periodically

[Timer]
OnCalendar={{ restic_postgres_schedule }}
Persistent=true

[Install]
WantedBy=timers.target
